# robot/pigpio.service
[Unit]
Description=Needed to run servo motor code

[Service]
Type=forking
ExecStart=pigpiod

[Install]
WantedBy=multi-user.target
